# DEBIAN/prerm
#!/bin/bash
/bin/systemctl stop trafficserver

# if 0 uninstall, if 1 upgrade
if [ "$1" = "0" ]; then
        /bin/systemctl disable trafficserver
fi
